# sat.cron
# run hourly the sat update

0 * * * *  root /usr/local/bin/sats.sh
