# sat.cron
# run hourly the sat update

1 * * * *  root /usr/local/bin/sats.sh
